# supervisor.conf
[program:nginx-app]
command = /usr/sbin/nginx

[program:mysqld]
command = /usr/sbin/mysqld --basedir=/usr --datadir=/var/lib/mysql --plugin-dir=/usr/lib/mysql/plugin --user=mysql --log-error=/var/log/mysql/error.log --pid-file=/var/run/mysqld/mysqld.pid --socket=/var/run/mysqld/mysqld.sock --port=3306

[program:web]
command=gunicorn -w 2 -b 127.0.0.1:9019 --worker-class=gevent izone.wsgi:application --log-file=/home/logs/web.log
directory = /home/izone/
redirect_stderr = true
stdout_logfile_maxbytes = 20MB
stdout_logfile_bakups = 10
stdout_logfile = /home/logs/web.log
autostart = true
autorestart = true